# local.d/50-foss-cloud-firstrun.start
../../usr/libexec/foss-cloud/foss-cloud-firstrun.sh