# cron.d/broker-restart
# MIN HOUR DAY MONTH DAYOFWEEK
# * * * * *
# Cronjob for the fc-broker
*/1 * * * * /usr/sbin/fc-restart-brokerd
00 01 * * * /etc/init.d/fc-brokerd restart

