# cron.d/cache-swap-clean
# MIN HOUR DAY MONTH DAYOFWEEK
# * * * * *
# Cronjob for cache and swap drop
0 */4 * * * /usr/sbin/fc-cache-swap-clean
