# cron.d/cache-drop
# MIN HOUR DAY MONTH DAYOFWEEK
# * * * * *
# Cronjob for cache drop
* */4 * * * /usr/sbin/fc-cache-drop
