# cron.d/check-updates
# MIN HOUR DAY MONTH DAYOFWEEK
# * * * * *
# Cronjob for the keeping /var/log clean
00 02 * * * /usr/sbin/fc-check-updates
