# cron.d/fc-log-maintenance
# MIN HOUR DAY MONTH DAYOFWEEK
# * * * * *
# Cronjob for the keeping /var/log clean
*/5 * * * * /usr/sbin/fc-log-backup
*/6 * * * * /usr/bin/db4.8_archive -a -h /var/lib/openldap-hdb/foss-cloud/ | xargs rm -v
