# cron.d/log-backup
# MIN HOUR DAY MONTH DAYOFWEEK
# * * * * *
# Cronjob for the fc-broker memory leaks
* */4 * * * /usr/sbin/fc-log-backup
