# cron.d/log-backup
# MIN HOUR DAY MONTH DAYOFWEEK
# * * * * *
# Cronjob for the keeping /var/log clean
0 */4 * * * /usr/sbin/fc-log-backup
