# cron.d/swap-clear
# MIN HOUR DAY MONTH DAYOFWEEK
# * * * * *
# Cronjob for swap clear
* */4 * * * /usr/sbin/fc-swap-clear
