# cron.d/utilization
# MIN HOUR DAY MONTH DAYOFWEEK
# * * * * *
# Cronjob for cache and swap drop and utilization
0 */4 * * * /usr/sbin/fc-utilization && /usr/sbin/pwfcck
