# init.d/fc-brokerd
#!/sbin/runscript
# Copyright 2012 FOSS-Group, Germany
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
    use net
    after slapd libvirtd
}

start() {
    local option_log

    [ -n "${LOGLEVEL}" ] && option_log="--log ${LOGLEVEL}"

    ebegin "Starting FOSS-Cloud Broker Daemon"
    start-stop-daemon --start \
        --user ${USER}:${GROUP} \
        --exec /usr/sbin/fc-brokerd \
        -- --config ${CONFIG} ${option_log} ${OPTIONS}
    eend $?
}

stop() {
    ebegin "Stopping FOSS-Cloud Broker Daemon"
    start-stop-daemon --stop \
        --exec /usr/sbin/fc-brokerd
    eend $?
}
